# aspects/systemd/files/.config/systemd/user/mailsync.timer
[Unit]
Description=Mail sync

[Timer]
OnBootSec=15m
OnUnitActiveSec=15m
Persistent=true

[Install]
WantedBy=timers.target
