# aspects/systemd/files/.config/systemd/user/mpris-proxy.service
[Unit]
Description=Forward Bluetooth media controls to MPRIS

[Service]
Type=simple
ExecStart=/usr/bin/mpris-proxy

[Install]
WantedBy=default.target
