# aspects/systemd/files/.config/systemd/user/newsup.timer
[Unit]
Description=Update newsboat

[Timer]
OnBootSec=1m
OnUnitActiveSec=15m
Persistent=true

[Install]
WantedBy=timers.target
