# aspects/systemd/files/.config/systemd/user/newsup.service
[Unit]
Description=Update newsboat

[Service]
ExecStart=/home/singularis/.local/bin/systemd/newsup
Environment="HOME=/home/singularis"
RestartSec=900

[Install]
WantedBy=default.target
