# aspects/autostart/xcompmgr.desktop
[Desktop Entry]
Type=Application
Exec=xcompmgr&
Terminal=false
StartupNotify=true
